# examples/ribeye/DEBIAN/prerm
#!/bin/bash
set -e

service ribeye stop
systemctl -q disable ribeye
systemctl -q daemon-reload
